# systemd/.config/systemd/user/bspwm-session.target
[Unit]
Description=Xorg Running
BindsTo=graphical-session.target
